# configs/nftables-exporter.service
[Unit]
Description=nftables exporter

[Service]
ExecStart=/usr/local/bin/nftables-exporter --group 0

[Install]
WantedBy=multi-user.target